# jobs/crontab
*/20  * * * *   root    shutdown -r now >> /home/pi/cron.log 2>&1
